# deploy/supervisord.conf
[program:wp-orchestrator-worker]
command=php /home/orchestrator/wp-orchestrator/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
directory=/home/orchestrator/wp-orchestrator
autostart=true
autorestart=true
user=orchestrator
numprocs=2
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/home/orchestrator/wp-orchestrator/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600
killasgroup=true
